# contrib/systemd/sist2-update.service
[Unit]
Description=sist2-update

[Service]
User=yatli
ExecStart=/bin/bash /usr/bin/sist2-update-all.sh
