# contrib/systemd/sist2-update.timer
[Unit]
Description=sist2-update

[Timer]
OnCalendar=*-*-* 3:00:00
Persistent=true
Unit=sist2-update.service

[Install]
WantedBy=timers.target
